# db-update.cron
0 8 * * * /home/appuser/.local/bin/cincanregistry list -r versions --silent >> /var/log/cron.log 2>&1
